# systemd/torch.socket
[Unit]
Description=Torch prometheus proxy socket

[Socket]
ListenStream=/run/torch.sock

[Install]
WantedBy=sockets.target
